# root/etc/sv/zram/finish
#!/bin/sh
set -e

DEV=$(cat /run/zram)
swapoff $DEV
zramctl --reset $DEV
rm /run/zram
