# root/etc/sv/duckdns/run
#!/bin/sh
set -u

while true
do
    curl -s "https://www.duckdns.org/update?domains=$DOMAIN&token=$TOKEN" >/dev/null
    sleep 300
done
