# root/etc/sv/duckdns/run
#!/bin/sh
set -u

while true
do
    IPV6="$(curl -s https://v6.ident.me)"
    [ -n "$IPV6" ] && curl -s "https://www.duckdns.org/update?domains=$DOMAIN&token=$TOKEN&ipv6=$IPV6"
    sleep 300
done
