# root/etc/sv/intel-undervolt/run
#!/bin/sh
intel-undervolt apply >/dev/null
exec pause
