# files/etc/profile.d/pipx.sh
export PIPX_GLOBAL_HOME=/usr/lib/pipx
export PIPX_GLOBAL_BIN_DIR=/usr/bin
export PIPX_MAN_DIR=/usr/share/man
